# tests/fixtures/host_etc/home/jdoe/.config/containers/systemd/dev-postgres.container
[Unit]
Description=Development PostgreSQL
After=network-online.target

[Container]
Image=docker.io/library/postgres:16
PublishPort=5432:5432
Volume=pgdata:/var/lib/postgresql/data:Z
Environment=POSTGRES_PASSWORD=devonly

[Service]
Restart=always

[Install]
WantedBy=default.target
