# tests/fixtures/host_etc/etc/containers/systemd/nginx.container
[Unit]
Description=nginx
[Container]
Image=quay.io/nginx
[Service]
