# tests/output/config/etc/systemd/system/cron-hourly-job.service
[Unit]
Description=Timer from cron etc/cron.d/hourly-job

[Service]
Type=oneshot
ExecStart=/bin/true
# FIXME: replace with actual command from cron
